# deploy/systemd/dxsb-planner.service
[Unit]
Description=DXSB Planner Sync and Scan
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/DXSB_lingonberry
ExecStart=/home/ubuntu/DXSB_lingonberry/scripts/planner_cycle.sh
